# root/etc/services.d/hubot/run
#!/usr/bin/with-contenv bash

cd /app || exit

exec /bin/ash /app/bin/hubot -a slack
